# supervisord.conf
[supervisord]

[program:amixer_web]
command=/usr/bin/python3 services/volume-control/run.py
autorestart=true
