# supervisord.conf
[supervisord]

[program:amixer_web]
command=/usr/bin/python3 amixer_web.py
autorestart=true
